# led_pekac/pekac.service
[Unit]
Description=LED wheel control

[Service]
Type=simple
ExecStart=/root/pekac.sh
Restart=always
